# r630_fan.service
[Unit]
Description=R630 fan control.

[Service]
ExecStart=/bin/bash /root/r630/r630_fan.sh

[Install]
WantedBy=multi-user.target

